# board/gt-be98/phase3/openrc/init.d/net-lan
#!/sbin/openrc-run
# === v8 ★REACHABILITY+DATAPATH★: wired LAN + Broadcom runner datapath + sshd ====
# v7 RESULT: open-init booted fully (OpenRC sysinit->boot->default), net-lan did
# config_switch + config_extwan + brctl addbr br0 + ip addr 10.0.0.8 (all [ok]),
# but NO traffic passed -> never reachable on slot2. On Broadcom, br0 + ip is NOT
# enough: the HW flow-accelerator/runner datapath (eth port <-> bridge <-> CPU)
# must be configured, which the ASUS rc START state does AFTER the bridge is up.
#
# v8 DATAPATH DIAGNOSIS (live br-0045 capture + rc source, sw_mode=3 AP/bridge):
#   The runner modules (pktrunner, bcm_enet, pktflow, cmdlist, rdpa*, bdmf) are
#   already loaded by bcm_boot_launcher (bcm-platform). The eth<->bridge<->CPU
#   datapath needs THREE steps beyond brctl/ip, all confirmed on the live device:
#   (a) ★ALLMULTI on br0 + every bridge member★ — rc start_lan brings the bridge
#       up via ifconfig(lan_ifname, IFUP|IFF_ALLMULTI, ...) and members are
#       ALLMULTI too (ifconfig br0/eth0 -> "UP BROADCAST RUNNING ALLMULTI"). The
#       bcm_enet+runner datapath floods/forwards to the CPU port via ALLMULTI;
#       v7 used `ip link set up` (NO ALLMULTI) -> the runner never bridged to CPU.
#   (b) ★fc enable★ — rc start_lan's fc_init() runs `fc enable` (flow-cache HW
#       accelerator) when !is_routing_enabled() (sw_mode=3 here). Live device:
#       `fcctl status` -> "HW Acceleration <Enabled>", fc_disable=0.
#   (c) ★rtpolicy auto ALL★ — rc init_main START state (init.c:26613, the
#       RTCONFIG_HND_ROUTER_BE_4916 step right after start_lan) runs
#       `rtpolicy auto ALL` -> applies the runner runtime policy from
#       /etc/rt_policy_info.d/. runner_disable=0.
#
# Open-path bring-up:
#   1. `rc config_switch` + `rc config_extwan` — standalone in-proc rc applets
#      (rc.c:4388/4394) that program the CLOSED Broadcom HW switch via the
#      ethswctl/vlanctl ABI. They do NOT signal PID1 -> work under OpenRC [V-src].
#   2. open start_lan stand-in: `brctl addbr br0` + add lan_ifnames members +
#      members UP with ALLMULTI + `ifconfig br0 ... IFF_ALLMULTI <ip> <nm>`, all
#      from nvram (lan_ifnames / lan_ipaddr / lan_netmask; nvram is up).
#   3. ★v8 datapath★: `fc enable` (flow cache) + `rtpolicy auto ALL` (runner) so
#      eth <-> br0 <-> CPU actually passes traffic -> SSH reachable.
#   4. sshd reachability: relaunch the :2223 RESCUE (dropbearmulti, /data hostkey;
#      replicates br-dropbear.sh — defensive, in case rc3.d ordering differs) and
#      start the merlin dropbear on :2222. Both key-auth, /data-persisted hostkey.
# Ordered AFTER bcm-platform (drivers + nvram loaded), BEFORE the rest.
description="GT-BE98 wired LAN + Broadcom runner datapath + sshd (open-init reachable)"

depend() {
	need nvram
	after bcm-platform
	before wifi-radio webui
	provide net
}

blog() { echo "$(cat /proc/uptime) net-lan: $*" >> /data/openrc-boot.log 2>/dev/null; }

# netmask dotted-quad -> CIDR prefix length (fallback 24).
nm2cidr() {
	local nm="$1" b n=0
	case "$nm" in *.*.*.*) ;; *) echo 24; return;; esac
	local IFS=.
	for b in $nm; do
		while [ "$b" -gt 0 ] 2>/dev/null; do n=$((n + (b & 1))); b=$((b >> 1)); done
	done
	[ "$n" -ge 1 ] && [ "$n" -le 32 ] && echo "$n" || echo 24
}

# Launch a dropbear listener on $1, binary $2, with a /data-persisted ed25519
# hostkey (stable fingerprint across boots/slots — same as br-dropbear.sh). A
# capped babysitter loop respawns it (OpenRC does not supervise it). Skip if the
# port is already bound (e.g. :2223 already up via rc3.d S28br-dropbear).
KEYDIR=/data/br/dropbear
start_dropbear() {
	local port="$1" dbin="$2" pidf="/tmp/dropbear-${port}.pid"
	[ -x "$dbin" ] || { blog "dropbear bin $dbin missing -> skip :$port"; return; }
	if netstat -ltn 2>/dev/null | grep -q ":${port} "; then
		blog "port :$port already listening -> not relaunching"
		return
	fi
	mkdir -p "$KEYDIR" 2>/dev/null
	[ -f "$KEYDIR/hostkey" ] || "$dbin" dropbearkey -t ed25519 -f "$KEYDIR/hostkey" >/dev/null 2>&1
	(
		N=0
		while [ "$N" -lt 20 ]; do
			"$dbin" dropbear -F -r "$KEYDIR/hostkey" -p "$port" -P "$pidf" >/dev/null 2>&1
			N=$((N + 1))
			sleep 5
		done
	) &
	blog "launched dropbear on :$port ($dbin, /data hostkey)"
}

start() {
	blog "START"

	# 1. HW switch + ext-WAN port programming (in-proc rc applets; graft-safe).
	ebegin "HW switch programming (rc config_switch)"
	rc config_switch >> /data/openrc-boot.log 2>&1
	eend $?
	ebegin "Ext-WAN port config (rc config_extwan)"
	rc config_extwan >> /data/openrc-boot.log 2>&1
	eend $?

	# 2. Open start_lan stand-in: br0 + lan_ifnames members + management IP.
	#    ★v8★: bridge AND every member brought up with ALLMULTI (IFF_ALLMULTI) —
	#    rc start_lan uses ifconfig(lan_ifname, IFUP|IFF_ALLMULTI, ...) and the
	#    members run ALLMULTI too (live br-0045: ifconfig br0/eth0 show ALLMULTI).
	#    The bcm_enet+runner datapath needs ALLMULTI to flood/forward to the CPU
	#    port; plain `ip link set up` (v7) left the runner with no CPU bridging.
	local br="br0" ifs ip nm cidr i
	ifs=$(nvram get lan_ifnames 2>/dev/null)
	ip=$(nvram get lan_ipaddr 2>/dev/null);  [ -n "$ip" ] || ip=10.0.0.8
	nm=$(nvram get lan_netmask 2>/dev/null); [ -n "$nm" ] || nm=255.255.255.0
	cidr=$(nm2cidr "$nm")
	blog "lan_ifnames='$ifs' lan_ipaddr=$ip lan_netmask=$nm (/$cidr)"
	ebegin "Bringing up $br ($ifs) -> $ip/$cidr (ALLMULTI)"
	if ! brctl show 2>/dev/null | grep -q "^${br}[[:space:]]"; then
		brctl addbr "$br" 2>/dev/null && blog "created bridge $br"
	fi
	for i in $ifs; do
		[ "$i" = "$br" ] && continue
		# member UP + ALLMULTI (ifconfig sets IFF_ALLMULTI; ip-link fallback)
		ifconfig "$i" up allmulti 2>/dev/null \
			|| { ip link set "$i" up 2>/dev/null; ip link set "$i" allmulticast on 2>/dev/null; }
		brctl addif "$br" "$i" 2>/dev/null && blog "added $i to $br (allmulti)"
	done
	# bring the bridge UP with the management IP AND ALLMULTI in one ifconfig —
	# this mirrors rc start_lan's ifconfig(lan_ifname, IFUP|IFF_ALLMULTI, ip, nm).
	ifconfig "$br" "$ip" netmask "$nm" up allmulti 2>/dev/null \
		|| { ip link set "$br" up 2>/dev/null
		     ip link set "$br" allmulticast on 2>/dev/null
		     ip addr add "$ip/$cidr" dev "$br" 2>/dev/null; }
	eend 0
	blog "bridge $br up (ALLMULTI) with $ip/$cidr"

	# 3. ★v8 Broadcom runner datapath★ — without these, br0+ip is logically up
	#    but the HW flow-accelerator never bridges eth<->br0<->CPU (the v7 gap).
	#    (a) fc enable: rc start_lan fc_init() -> `fc enable` (flow cache HW accel,
	#        AP/bridge mode, fc_disable=0). `fc` is /bin/fc -> fcctl on the device.
	#    (b) rtpolicy auto ALL: rc init_main START state (init.c:26613) -> applies
	#        the runner runtime policy from /etc/rt_policy_info.d/.
	ebegin "Broadcom datapath: fc enable + rtpolicy auto ALL"
	if [ "$(nvram get fc_disable 2>/dev/null)" != "1" ]; then
		fc enable >> /data/openrc-boot.log 2>&1 && blog "fc enable (flow cache HW accel)"
	else
		blog "fc_disable=1 -> skip fc enable"
	fi
	if [ "$(nvram get runner_disable 2>/dev/null)" != "1" ]; then
		rtpolicy auto ALL >> /data/openrc-boot.log 2>&1 && blog "rtpolicy auto ALL (runner policy)"
	else
		blog "runner_disable=1 -> skip rtpolicy"
	fi
	eend 0
	blog "datapath configured — eth <-> $br <-> CPU should now pass traffic"

	# 4. sshd reachability: :2223 rescue (defensive; rc3.d may already have it) + :2222.
	ebegin "Starting sshd (:2223 rescue + :2222)"
	start_dropbear 2223 /usr/br/sbin/dropbearmulti
	start_dropbear 2222 /usr/sbin/dropbear
	eend 0

	blog "DONE — datapath up; open-init reachable on :2222 (+ :2223 rescue) at $ip"

	# 5. ★v9 SELF-DIAGNOSTICS★ — v8 ran every datapath step + dropbear "listening"
	#    yet NO traffic passed on :2222/:2223, and live debug is impossible (the
	#    image runs unattended on the trial, then the dead-man reverts). So the
	#    image must SELF-REPORT: dump the full network + Broadcom-datapath state
	#    and a connectivity self-test to PERSISTENT /data/net-diag.log (survives
	#    the revert). Backgrounded + every probe time-bounded so it can NEVER hang
	#    the boot. The next trial's /data/net-diag.log then says whether the
	#    datapath passes traffic (gateway ping) and whether sshd actually listens.
	run_diagnostics "$br" "$ip" &

	return 0
}

# === v9 self-diagnostics ====================================================
# Append the live network + datapath + listener + connectivity state to a
# PERSISTENT log. PRIVACY: never log keys/MACs — every dump is piped through
# redact_mac (xx:xx:xx:xx:xx:xx) and we never touch the dropbear hostkey/
# fingerprint. Each command is guarded `2>&1`, skipped if its binary is absent,
# and time-bounded (run_to) so a stuck tool cannot wedge the diagnostics.
DIAG=/data/net-diag.log

# redact every MAC (and bare 6-octet hwaddr) so net-diag.log carries no MACs.
redact_mac() {
	sed -E 's/([0-9a-fA-F]{2}:){5}[0-9a-fA-F]{2}/xx:xx:xx:xx:xx:xx/g'
}

# run_to SECS CMD... — run a command bounded by `timeout` (busybox has it);
# falls back to a background kill if timeout is unavailable. Output is appended
# (already redacted by the caller's pipeline).
run_to() {
	local s="$1"; shift
	if command -v timeout >/dev/null 2>&1; then
		timeout "$s" "$@" 2>&1
	else
		( "$@" 2>&1 ) & local p=$!
		( sleep "$s"; kill "$p" 2>/dev/null ) & local k=$!
		wait "$p" 2>/dev/null; kill "$k" 2>/dev/null
	fi
}

# diag SECS "LABEL" CMD... — emit a labeled section iff CMD's binary exists,
# redacting MACs. Skips (with a note) if the binary is missing.
diag() {
	local s="$1" label="$2"; shift 2
	{
		echo "----- $label -----"
		if command -v "$1" >/dev/null 2>&1; then
			run_to "$s" "$@" | redact_mac
		else
			echo "(skip: $1 not present)"
		fi
		echo
	} >> "$DIAG" 2>/dev/null
}

# ping_test LABEL TARGET — ping with a clear PASS/FAIL line (the key datapath
# probe). 3 echo requests, 2s per-packet wait, whole thing bounded to 12s.
ping_test() {
	local label="$1" tgt="$2" out rc
	{
		echo "----- PING $label ($tgt) -----"
		if [ -z "$tgt" ]; then
			echo "RESULT: SKIP $label — no target"
			echo; return
		fi
		if ! command -v ping >/dev/null 2>&1; then
			echo "RESULT: SKIP $label — ping not present"
			echo; return
		fi
		out=$(run_to 12 ping -c3 -W2 "$tgt"); rc=$?
		echo "$out" | redact_mac
		if [ "$rc" -eq 0 ] && echo "$out" | grep -q ' 0% packet loss'; then
			echo "RESULT: PASS $label ($tgt) reachable"
		else
			echo "RESULT: FAIL $label ($tgt) UNREACHABLE (rc=$rc)"
		fi
		echo
	} >> "$DIAG" 2>/dev/null
}

run_diagnostics() {
	local br="$1" ip="$2" gw
	gw=$(nvram get lan_gateway 2>/dev/null)
	{
		echo "================================================================"
		echo "GT-BE98 open-init-v9 net-lan SELF-DIAGNOSTICS"
		echo "uptime: $(cat /proc/uptime 2>/dev/null)"
		echo "date:   $(date 2>/dev/null)"
		echo "br=$br ip=$ip lan_gateway=${gw:-<unset>}"
		echo "================================================================"
	} >> "$DIAG" 2>/dev/null

	# --- L2/L3 network state (MACs redacted) ---
	diag 5  "ip -d addr show"        ip -d addr show
	diag 5  "ip route show"          ip route show
	diag 5  "ip neigh"               ip neigh
	diag 5  "arp -an"                arp -an
	diag 5  "brctl show"             brctl show
	diag 5  "brctl showmacs $br"     brctl showmacs "$br"
	diag 5  "/proc/net/dev"          cat /proc/net/dev

	# --- Broadcom switch / runner / flow-accelerator datapath state ---
	# ethswctl/ethctl: link + port state for eth0 (the LAN uplink to the switch).
	diag 5  "ethswctl mode (getifname unit0 port0)" ethswctl getifname 0 0
	diag 5  "ethctl eth0 media-type"  ethctl eth0 media-type
	diag 5  "ethctl eth0 phy-power"   ethctl eth0 phy-power
	diag 5  "vlanctl --dump-table"    vlanctl --dump-table
	diag 5  "fcctl status"            fcctl status
	# rtpolicy has no 'show' subcmd; the applied runner policy lives in this dir.
	diag 5  "rtpolicy config dir (/etc/rt_policy_info.d)" ls -la /etc/rt_policy_info.d
	# loaded runner/flow modules (datapath drivers) — confirms the HW path exists.
	diag 5  "loaded runner/flow modules" sh -c 'grep -E "pktrunner|pktflow|bcm_enet|cmdlist|rdpa|bdmf|archer" /proc/modules'
	diag 5  "/proc/driver/enet listing" ls -la /proc/driver/enet

	# --- ★listening sockets: confirm dropbear is REALLY listening :2222+:2223★ ---
	{
		echo "----- LISTENING TCP SOCKETS (netstat -ltn || ss -ltn) -----"
		if command -v netstat >/dev/null 2>&1; then
			run_to 5 netstat -ltn | redact_mac
		elif command -v ss >/dev/null 2>&1; then
			run_to 5 ss -ltn | redact_mac
		else
			echo "(skip: neither netstat nor ss present)"
		fi
		echo
		# explicit per-port verdict for the two dropbear listeners
		for p in 2222 2223; do
			if { netstat -ltn 2>/dev/null || ss -ltn 2>/dev/null; } | grep -q ":${p} "; then
				echo "RESULT: PASS sshd LISTENING on :$p"
			else
				echo "RESULT: FAIL sshd NOT listening on :$p"
			fi
		done
		echo
	} >> "$DIAG" 2>/dev/null

	# --- ★CONNECTIVITY SELF-TEST (the key datapath probe)★ ---
	ping_test "loopback"        127.0.0.1
	ping_test "gateway-10.0.0.1" 10.0.0.1
	ping_test "nvram-lan_gateway" "$gw"

	echo "===== net-diag complete =====" >> "$DIAG" 2>/dev/null
	blog "self-diagnostics written to $DIAG"
}
